# supervisord.d/huskar_console.ini
[program:huskar_console]
command=/opt/huskar_console/bin/entrypoint.sh
autorestart=true
startretries=3
stopsignal=TERM
directory=/var/lib/nginx
user=root
group=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
stopasgroup=true
killasgroup=true
