# supervisord.d/mariadb.ini
[program:mariadb]
command=/usr/libexec/mysqld --basedir=/usr
autorestart=true
startretries=3
stopsignal=TERM
directory=/var/lib/mysql
user=mysql
group=mysql
stdout_logfile=/var/log/mariadb/supervisor.log
stdout_logfile_maxbytes=3MB
stdout_logfile_backups=5
redirect_stderr=true
stopasgroup=true
killasgroup=true
