# supervisord.d/zookeeper.ini
[program:zookeeper]
command=/usr/bin/zookeeper-server start-foreground
autorestart=true
startretries=3
stopsignal=TERM
directory=/var/lib/zookeeper
user=zookeeper
group=zookeeper
stdout_logfile=/var/log/zookeeper/supervisor.log
stdout_logfile_maxbytes=3MB
stdout_logfile_backups=5
redirect_stderr=true
stopasgroup=true
killasgroup=true
